# shanfomi_fable/v_20260618_232334/src/launchd/ai.bittlebits.fomi4me.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.bittlebits.fomi4me.digest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>src/digest_builder.py</string>
        <string>--send</string>
    </array>
    <!-- Set WorkingDirectory to the repo root when installing (make install-launchd). -->
    <key>WorkingDirectory</key>
    <string>__REPO__</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__REPO__/logs/digest.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/logs/digest.log</string>
</dict>
</plist>
